# modules/jinteki/templates/systemd/jinteki-site.service.epp
[Unit]
Description=Site server for Jinteki
After=network.target

[Service]
Environment="NODE_ENV=production"
Type=simple
ExecStart=/usr/bin/coffee server.coffee
PrivateDevices=yes
User=jinteki
WorkingDirectory=<%= "${jinteki::repodir}" %>

[Install]
WantedBy=multi-user.target

